# .agents/daemons/macos-host/com.ostwin.macos-host.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ostwin.macos-host</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>OSTWIN_HOME/.agents/daemons/macos-host/host-daemon.sh</string>
  </array>

  <!-- Run immediately on login and keep alive if it exits -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Log files (OSTWIN_HOME is substituted by install.sh) -->
  <key>StandardOutPath</key>
  <string>OSTWIN_HOME/logs/macos-host.log</string>
  <key>StandardErrorPath</key>
  <string>OSTWIN_HOME/logs/macos-host-err.log</string>

  <!-- Environment available to the daemon process -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>OSTWIN_HOME</key>
    <string>OSTWIN_HOME</string>
    <key>OSTWIN_SOCKET</key>
    <string>/tmp/ostwin-macos-host.sock</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
